# deploy/bracer-bot/bracer-bot.service
[Unit]
Description=Bracer Support Bot
After=network.target snap.rocketchat-server.rocketchat-server.service

[Service]
Type=simple
User=bracer-bot
Group=bracer-bot
WorkingDirectory=/opt/bracer-bot
EnvironmentFile=/opt/bracer-bot/.token-env
EnvironmentFile=/opt/bracer-bot/.env.plain
ExecStart=op run --env-file=/opt/bracer-bot/.env.op -- /opt/bracer-bot/venv/bin/uvicorn bot:app --host 127.0.0.1 --port 5051 --no-access-log
Restart=always
RestartSec=5
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
